# ccm-outreach.timer
[Unit]
Description=Run CCM Outreach Engine every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target